# infra/systemd/shipsquares-updater.service
[Unit]
Description=ShipSquares self-updater (download, migrate, swap, restart, roll back)
# Triggered by shipsquares-updater.path. Runs as root, outside the control-plane
# process, so it can swap /opt/shipsquares/current and restart the service.

[Service]
Type=oneshot
ExecStart=/opt/shipsquares/bin/update.sh
# The updater must write /opt/shipsquares, run migrations as the service user via
# `sudo -u shipsquares`, and drive `systemctl restart` — so the aggressive
# confinements are deliberately OMITTED (they would break it): NoNewPrivileges and
# RestrictSUIDSGID block sudo's setuid transition, and ProtectSystem=strict would
# block the /opt swap. What IS safe to lock down, we do — the updater never touches
# home dirs, kernel tunables/modules, cgroups, the clock, hostname, or realtime:
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
LockPersonality=true
